# wifi-monitor.openrc
#!/sbin/openrc-run

pidfile="/run/wifi-monitor.pid"
command="/usr/sbin/wifi-monitor"
command_args="-d"

